# templates/rc.d.in
#!/bin/sh

# PROVIDE: minecraft
# REQUIRE: LOGIN
# KEYWORD: shutdown

. /etc/rc.subr

name="minecraft"
rcvar=minecraft_enable

load_rc_config $name

: ${minecraft_enable:="NO"}
: ${minecraft_user:="@MINECRAFT_USER@"}
: ${minecraft_group:="@MINECRAFT_GROUP@"}
: ${minecraft_dir:="@MINECRAFT_DIR@"}
: ${service_script:="@SERVICE_SCRIPT@"}

start_cmd="$service_script start"
stop_cmd="$service_script stop"
status_cmd="$service_script status"
log_cmd="$service_script log"
attach_cmd="$service_script attach"
cmd_cmd="$service_script cmd"
reload_cmd="$service_script reload"
extra_commands="log attach cmd reload"

run_rc_command "$@"
